# fedora/systemd/.config/systemd/user/app-dev.lizardbyte.app.Sunshine.service.d/override.conf
# Drop-in for the rpm's Sunshine unit. The gamescope stream session starts this
# unit (via steam-session, GS_SUNSHINE=1), but the per-user systemd environment
# still carries WAYLAND_DISPLAY=wayland-1 from the Sway session that last
# imported it. Sunshine then auto-selects wlr (Wayland) capture and fails:
#   Error: [wayland] Couldn't connect to Wayland display: wayland-1
# -> Fatal: Unable to find display or encoder during startup.
#
# gamescope owns the DRM scanout directly and doesn't expose wlr-screencopy, so
# the correct path is KMS capture (the cap_sys_admin file-cap on the binary).
# Force capture=kms and clear the stale Wayland/X vars so KMS enumerates the
# real DRM connector. KMS capture also works fine from the Sway desktop, so this
# drop-in is harmless to the occasional manual run there.
[Service]
UnsetEnvironment=WAYLAND_DISPLAY DISPLAY
ExecStart=
ExecStart=/usr/bin/sunshine capture=kms
