# fedora/config/openrgb/rgb.service
[Unit]
Description=Set OpenRGB device 0 lighting at boot
# System service for early-boot RGB (runs before any user logs in), the
# OpenRGB analog of lactd. Installed by setup-openrgb.sh; not stowed since
# system units live under /etc. Edit --color to taste; 000000 turns it off.
After=multi-user.target

[Service]
Type=oneshot
# OpenRGB writes a ./logs dir relative to CWD; the default CWD (/) is read-only
# on Atomic, so point it at a writable runtime dir.
WorkingDirectory=/tmp
ExecStart=/usr/bin/openrgb --device 0 --color 000000
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
