# fedora/systemd/.config/systemd/user/niri-session.target
[Unit]
Description=Niri Session Services
